# cron/crontab
*/5 * * * * sh ~/effect_scanner/cron/scan.sh > ~/effect_scanner/cron.log 2>&1
