# cron/crontab
* * * * * sh ~/effect_scanner/cron/scan.sh > ~/effect_scanner/scan_simple.log 2>&1
